# peerjs-server.service
[Unit]
Description=PeerJS Server for CBT (Production)
After=network.target

[Service]
Type=simple
User=procbt
WorkingDirectory=/var/www/html/cbt-new
ExecStart=/usr/bin/node /var/www/html/cbt-new/peerjs-server.js
Restart=always
RestartSec=10
Environment="PORT=9000"
Environment="NODE_ENV=production"

[Install]
WantedBy=multi-user.target
